# deploy/headroom.service.template
[Unit]
Description=Headroom proxy (context compression layer for AI agents)
Documentation=https://github.com/chopratejas/headroom
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# __HEADROOM_BIN__ 与 __PORT__ 由 deploy-linux.sh 在安装时替换。
ExecStart=__HEADROOM_BIN__ proxy --host 0.0.0.0 --port __PORT__
Environment=HEADROOM_HOST=0.0.0.0
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
